# images/server1/bootstrap.service
[Unit]
Description=Set up system

[Service]
Type=oneshot
ExecStart=/usr/bin/bootstrap

[Install]
WantedBy=multi-user.target
